# bin/cron.hourly/stream_refresh-admins.sh
#!/bin/bash

cd /var/www/sparkrebel.com/current

app/console stream:refresh "5033bb99a32871764f000000" user # dabrowski85@gmail.com
app/console stream:refresh "4efadbbe411e20d64c00000f" user # liat.ts@gmail.com
app/console stream:refresh "4efafe9b411e20dc4c000018" user # leeked@gmail.com
app/console stream:refresh "4efb509c411e20df4c00000e" user # elad@eladbaron.com
app/console stream:refresh "4efcd2d7411e206e65000000" user # lindsay.anvik@gmail.com
app/console stream:refresh "4f032383411e20460e000002" user # plumwillow@ad7six.com
app/console stream:refresh "4f16e26f411e206b4400001e" user # omerfb@plumwillow.com
app/console stream:refresh "4f4ec4c5c56da10e5e000000" user # saddor@gmail.com
app/console stream:refresh "4f73356fc56da15d06000012" user # dshadow@zort.net
app/console stream:refresh "4f8f0ec8c56da11951000000" user # dabrowski85@gmail.com
app/console stream:refresh "501808f3140ba0f948000001" user # shulamitkashy@gmail.com

app/console stream:refresh "5033bb99a32871764f000000" brand # dabrowski85@gmail.com
app/console stream:refresh "4efadbbe411e20d64c00000f" brand # liat.ts@gmail.com
app/console stream:refresh "4efafe9b411e20dc4c000018" brand # leeked@gmail.com
app/console stream:refresh "4efb509c411e20df4c00000e" brand # elad@eladbaron.com
app/console stream:refresh "4efcd2d7411e206e65000000" brand # lindsay.anvik@gmail.com
app/console stream:refresh "4f032383411e20460e000002" brand # plumwillow@ad7six.com
app/console stream:refresh "4f16e26f411e206b4400001e" brand # omerfb@plumwillow.com
app/console stream:refresh "4f4ec4c5c56da10e5e000000" brand # saddor@gmail.com
app/console stream:refresh "4f73356fc56da15d06000012" brand # dshadow@zort.net
app/console stream:refresh "4f8f0ec8c56da11951000000" brand # dabrowski85@gmail.com
app/console stream:refresh "501808f3140ba0f948000001" brand # shulamitkashy@gmail.com
